# .config/systemd/user/check-battery.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/chris/.config/systemd/user/check-battery.sh
